# supervisor.conf
[program:httpstatus]
command=php console.php /controllers/internals/WebsiteStatus init
directory=/var/www/html/httpstatus/
autostart=true
autorestart=true